# .backup/coauthor-backup.service
[Unit]
Description=Run coauthor's backup once

[Service]
Type=oneshot
ExecStart=/usr/bin/sh %h/dev/coauthor/.backup/backup.sh
