# .backup/coauthor-vacuum.service
[Unit]
Description=Run coauthor's vacuum once

[Service]
Type=oneshot
ExecStart=/usr/bin/bash %h/dev/coauthor/.backup/vacuum.sh
